# s6-config/etc/services.d/openvpn/finish
#!/command/with-contenv bash
# If service exited with error, stop the container
if [ "$1" -ne 0 ]; then
    echo "openvpn exited with code $1, stopping container..."
    /run/s6/basedir/bin/halt
fi
