# s6-config/etc/services.d/openvpn/finish
#!/command/execlineb -S1
# Exit code 256 tells s6 not to restart this service
ifelse { s6-test ${1} -ne 0 } { s6-svscanctl -t /var/run/s6/services }
