# s6-config/etc/services.d/openvpn/run
#!/command/with-contenv bash
exec openvpn --config /etc/openvpn/client.ovpn
